# systemd/check-disk-usage.service
[Unit]
Description=Check disk usage and warn if over threshold
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash -lc 'TH=80; for p in /data /var/lib/marketbrewer; do if [ -d "$p" ]; then used=$(df -P "$p" | awk "NR==2{gsub(\"%\", \"\", $5); print $5}"); if [ "$used" -ge "$TH" ]; then echo "WARNING: $p usage ${used}% exceeds ${TH}%"; fi; fi; done'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target