# systemd/seo-worker.service
[Unit]
Description=MarketBrewer SEO Platform Content Generation Worker
Documentation=https://github.com/dmv7zero3/marketbrewer-seo-platform
After=network.target seo-api.service ollama.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/marketbrewer-seo-platform
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NODE_ENV=production"
Environment="API_URL=http://localhost:3001"
Environment="OLLAMA_BASE_URL=http://localhost:11434"
Environment="OLLAMA_MODEL=llama3.2"
Environment="WORKER_POLL_INTERVAL=5000"
Environment="WORKER_BATCH_SIZE=2"
Environment="LOG_LEVEL=warn"
EnvironmentFile=/var/www/marketbrewer-seo-platform/.env

ExecStart=/usr/bin/node packages/worker/dist/index.js

# Restart policy
Restart=on-failure
RestartSec=15

# Resource limits
LimitNOFILE=8192
LimitNPROC=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seo-worker

# Security
PrivateTmp=yes
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
